# src/supervisord.conf
[supervisord]
nodaemon=false

[program:patcher]
priority=30
directory=/server
command=ruby main.rb
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/patcher.log
stderr_logfile=/var/log/patcher.err

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix://var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
